# cvs.conf
[program:cvs]
directory=/var/python/cvs
command=/var/python/cvs/cvs/bin/python compteur.py
user=david
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true